# restic/restic-backup.timer
[Unit]
Description=Daily Restic Backup to MinIO

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
