# networking/unbound/systemd/update-unbound-hints.timer
[Unit]
Description=Run Unbound Root Hints Update Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
